# systemd/breaktime.service
[Unit]
Description=Breaktime - Automated break scheduling system
Documentation=https://github.com/benjaminpeeters/breaktime
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=SCRIPT_PATH --daemon
Restart=always
RestartSec=10
TimeoutStopSec=30

# Environment
Environment=HOME=%h
Environment=XDG_RUNTIME_DIR=/run/user/%i
Environment=DISPLAY=:0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=breaktime

# Security settings (relaxed for user service)
NoNewPrivileges=true
PrivateTmp=true

# Allow access to user configuration and crontab
ReadWritePaths=%h/.config/breaktime

[Install]
WantedBy=default.target